# test123/rootfs/etc/services.d/test123/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the test123 service - GPS data echo app
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Declare variables
declare gps_device baud_rate

# Get configuration options
gps_device=$(bashio::config 'gps_device')
baud_rate=$(bashio::config 'baud_rate')

# Default values
gps_device="${gps_device:-/dev/ttyACM0}"
baud_rate="${baud_rate:-9600}"

bashio::log.info "Starting GPS data echo app..."
bashio::log.info "GPS device: ${gps_device}"
bashio::log.info "Baud rate: ${baud_rate}"

# Run the GPS echo program
exec /usr/bin/gps_echo
